# etc/cron.d/ingest_xdmod
# Run the hourly jobs
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
10 * * * * root /orcacluster/xdmod/ingest
