# etc/cron.d/xdmod
# Run the hourly xdmod collect
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=mewan@pdx.edu
00 * * * * root /cluster/xdmod/collect
